# raspsuraler.service
[Unit]
Description=RaspSurAlert — Surveillance aérienne locale
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=david
WorkingDirectory=/home/david/raspsuraler
ExecStart=/usr/bin/python3 /home/david/raspsuraler/server.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
